# Library/LaunchAgents/com.koekeishiya.skhd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.koekeishiya.skhd</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/skhd</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/tgelin01/.asdf/shims:/opt/homebrew/Cellar/asdf/0.12.0/libexec/bin:/usr/local/bin:/Users/tgelin01/.asdf/installs/flutter/3.10.5-stable/bin:/Users/tgelin01/.maestro/bin:/opt/homebrew/opt/llvm/bin:/Users/tgelin01/.local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/tgelin01/.asdf/installs/rust/1.70.0/bin:/Users/tgelin01/Library/Android/sdk/platform-tools:/Users/tgelin01/Library/Android/sdk/tools/bin:/Users/tgelin01/Library/Android/sdk/emulator:/Users/tgelin01/Library/Android/sdk/tools:/Users/tgelin01/.jenv/bin:/Users/tgelin01/.local/bin:/opt/homebrew/bin:/usr/sbin:/usr/local/bin:/System/Cryptexes/App/usr/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/local/bin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/bin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/appleinternal/bin:/Library/Apple/usr/bin:/opt/homebrew/bin :/Applications/kitty.app/Contents/MacOS</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
 	     <false/>
 	     <key>Crashed</key>
 	     <true/>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/skhd_tgelin01.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/skhd_tgelin01.err.log</string>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>Nice</key>
    <integer>-20</integer>
</dict>
</plist>