# Library/LaunchAgents/com.github.facebook.watchman.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.github.facebook.watchman</string>
    <key>Disabled</key>
    <false/>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/watchman</string>
        <string>--foreground</string>
        <string>--logfile=/opt/homebrew/var/run/watchman/tgelin01-state/log</string>
        <string>--log-level=1</string>
        <string>--sockname=/opt/homebrew/var/run/watchman/tgelin01-state/sock</string>
        <string>--statefile=/opt/homebrew/var/run/watchman/tgelin01-state/state</string>
        <string>--pidfile=/opt/homebrew/var/run/watchman/tgelin01-state/pid</string>
    </array>
    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/var/folders/9b/cg7wd0fn1jq7cfs1q_4sd90c0000gq/T/yarn--1692894351127-0.5779045875071254:/Users/tgelin01/Projects/enterprise-react-native/node_modules/.bin:/Users/tgelin01/.config/yarn/link/node_modules/.bin:/Users/tgelin01/.yarn/bin:/Users/tgelin01/.asdf/installs/nodejs/18.17.1/libexec/lib/node_modules/npm/bin/node-gyp-bin:/Users/tgelin01/.asdf/installs/nodejs/18.17.1/lib/node_modules/npm/bin/node-gyp-bin:/Users/tgelin01/.asdf/installs/nodejs/18.17.1/bin/node_modules/npm/bin/node-gyp-bin:/Users/tgelin01/.asdf/plugins/nodejs/shims:/Users/tgelin01/.asdf/installs/nodejs/18.17.1/bin:/Users/tgelin01/.asdf/installs/flutter/3.10.5-stable/bin:/Users/tgelin01/.maestro/bin:/Users/tgelin01/.local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/opt/llvm/bin:/Users/tgelin01/Library/Android/sdk/platform-tools:/Users/tgelin01/Library/Android/sdk/tools/bin:/Users/tgelin01/Library/Android/sdk/emulator:/Users/tgelin01/Library/Android/sdk/tools:/Users/tgelin01/.jenv/bin:/Users/tgelin01/.local/bin:/opt/homebrew/bin:/usr/sbin:/usr/local/bin:/Users/tgelin01/.asdf/shims:/opt/homebrew/Cellar/asdf/0.12.0/libexec/bin:/usr/local/bin:/Users/tgelin01/.asdf/installs/flutter/3.10.5-stable/bin:/Users/tgelin01/.maestro/bin:/Users/tgelin01/.local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/opt/llvm/bin:/Users/tgelin01/Library/Android/sdk/platform-tools:/Users/tgelin01/Library/Android/sdk/tools/bin:/Users/tgelin01/Library/Android/sdk/emulator:/Users/tgelin01/Library/Android/sdk/tools:/Users/tgelin01/.jenv/bin:/Users/tgelin01/.local/bin:/opt/homebrew/bin:/usr/sbin:/usr/local/bin:/System/Cryptexes/App/usr/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/local/bin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/bin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/appleinternal/bin:/Library/Apple/usr/bin:/opt/homebrew/bin :/Applications/kitty.app/Contents/MacOS</string>
    </dict>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>Nice</key>
    <integer>-5</integer>
</dict>
</plist>
