# Library/LaunchAgents/com.github.facebook.watchman.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.github.facebook.watchman</string>
    <key>Disabled</key>
    <false/>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/watchman</string>
        <string>--foreground</string>
        <string>--logfile=/opt/homebrew/var/run/watchman/thalesgelinger-state/log</string>
        <string>--log-level=1</string>
        <string>--sockname=/opt/homebrew/var/run/watchman/thalesgelinger-state/sock</string>
        <string>--statefile=/opt/homebrew/var/run/watchman/thalesgelinger-state/state</string>
        <string>--pidfile=/opt/homebrew/var/run/watchman/thalesgelinger-state/pid</string>
    </array>
    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/var/folders/mx/cg5v5hxn4sg9ppydpy37sv7m0000gp/T/yarn--1719685916018-0.6056764914464372:/Users/thalesgelinger/Projects/ProjetoTeste/node_modules/.bin:/Users/thalesgelinger/.config/yarn/link/node_modules/.bin:/Users/thalesgelinger/.yarn/bin:/Users/thalesgelinger/.local/share/mise/installs/node/20.14.0/libexec/lib/node_modules/npm/bin/node-gyp-bin:/Users/thalesgelinger/.local/share/mise/installs/node/20.14.0/lib/node_modules/npm/bin/node-gyp-bin:/Users/thalesgelinger/.local/share/mise/installs/node/20.14.0/bin/node_modules/npm/bin/node-gyp-bin:/Users/thalesgelinger/.local/share/mise/installs/node/lts/bin:/Users/thalesgelinger/.local/share/mise/installs/rust/latest/bin:/Users/thalesgelinger/.local/share/mise/installs/go/1.22.4/bin:/Users/thalesgelinger/.local/share/mise/installs/go/1.22.4/go/bin:/Users/thalesgelinger/.local/share/mise/installs/neovim/latest/bin:/Users/thalesgelinger/.local/share/mise/installs/flutter/latest/bin:/Users/thalesgelinger/.local/share/mise/installs/java/17/bin:/Users/thalesgelinger/.local/share/mise/installs/cocoapods/latest/bin:/Users/thalesgelinger/.local/share/mise/installs/bundler/latest/exe:/Users/thalesgelinger/Library/Android/sdk/ndk/26.1.10909125/toolchains/llvm/prebuilt/darwin-x86_64/bin:/Users/thalesgelinger/.local/share/mise/shims:/Users/thalesgelinger/.maestro/bin:/Users/thalesgelinger/.maestro/bin:/opt/homebrew/opt/llvm/bin:/Users/thalesgelinger/.local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/thalesgelinger/Library/Android/sdk/platform-tools:/Users/thalesgelinger/Library/Android/sdk/tools/bin:/Users/thalesgelinger/Library/Android/sdk/emulator:/Users/thalesgelinger/Library/Android/sdk/tools:/Users/thalesgelinger/.jenv/bin:/Users/thalesgelinger/.local/bin:/opt/homebrew/bin:/usr/sbin:/usr/local/bin:/opt/homebrew/opt/libpq/bin:/usr/local/bin:/System/Cryptexes/App/usr/bin:/usr/bin:/bin:/usr/sbin:/sbin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/local/bin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/bin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/appleinternal/bin:/Library/Apple/usr/bin:/opt/homebrew/bin:/Applications/kitty.app/Contents/MacOS</string>
    </dict>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>Nice</key>
    <integer>-5</integer>
</dict>
</plist>
