# Library/LaunchAgents/timesheet.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>timesheet.plist</string>
    <key>Program</key>
    <string>~/.local/bin/timesheet</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>55</integer>
        <key>Weekday</key>
        <integer>4</integer> 
    </dict>
</dict>
</plist>
